# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-rclone/finish
#!/bin/sh
# 卸载挂载点
fusermount -uz "/n3_storage/webdav_alist" >/dev/null 2>&1